# levels/app.conf
[program:portal]
directory=/opt/portal
command=/opt/portal/portal/bin/python app.py
